# devudo/templates/default/hosting_queue_runner.conf.erb
[program:hosting_queue_runner]
; Adjust the next line to point to where you copied the script.
command=<%= node[:aegir][:hosting_queue_runner_path] %>
user=aegir
numprocs=1
stdout_logfile=/var/log/hosting_queue_runner
autostart=TRUE
autorestart=TRUE
; Tweak the next line to match your environment.
environment=HOME="<%= node[:aegir][:dir] %>",USER="aegir",DRUSH_COMMAND="/usr/share/drush/drush"